# utils/cursor/update-cursor.service
[Unit]
Description=Update Cursor

[Service]
ExecStart=/opt/cursor/update-cursor.sh
Type=oneshot

[Install]
WantedBy=multi-user.target
